# file_sync/rootfs/etc/services.d/ssl_sync/finish
#!/bin/bash

echo "[ssl_sync] service stopped (code: $1, signal: $2)"
