# file_sync/rootfs/etc/services.d/ssl_sync/finish
#!/bin/bash
bashio::log.info "Service stopped (code: $1, signal: $2)"
exit 0
