# file_sync/rootfs/etc/services.d/ssl_sync/run
#!/bin/bash
set -e

echo "[ssl_sync] service starting..."

# Запуск основного скрипта
exec /root/run.sh

