# file_sync/rootfs/etc/services.d/ssl_sync/run
#!/bin/bash
bashio::log.info "[ssl_sync] Stopped (code: $1, signal: $2)"
exit 0
