# file_sync/rootfs/etc/services.d/ssl_sync/run
#!/usr/bin/with-contenv bash
# S6: exec daemon (простой, без set -e — s6 handles)
exec 2>&1  # Stdout to logs
exec /usr/local/bin/run.sh
