# file_sync/rootfs/etc/services.d/ssl_sync/run
#!/usr/bin/with-contenv bash
exec 2>&1
exec /usr/local/bin/run.sh  # <-- Здесь /usr/local/bin/, не /root/!
