# file_sync/rootfs/etc/services.d/file_sync/finish
#!/usr/bin/with-contenv bash
# finish called by s6 on service stop
echo "File Sync service stopped (code: $1, signal: $2)"
exit 0
