# file_sync/rootfs/etc/services.d/file_sync/finish
#!/bin/bash
bashio::log.info "File Sync service stopped (code: $1, signal: $2)"
exit 0
