# file_sync/rootfs/etc/services.d/file_sync/run
#!/command/with-contenv bash
# Обратите внимание: /command/, а не /usr/bin/

exec /run.sh
