# file_sync/rootfs/etc/services.d/file_sync/run
#!/usr/bin/with-contenv bash
# Этот скрипт запускает s6-overlay, он ДОЛЖЕН быть максимально простым

exec /run.sh
