# file_sync/rootfs/etc/services.d/file_sync/run
#!/usr/bin/with-contenv bash
# Фикс s6: simple exec, no set -e, exec 2>&1 для логов
exec 2>&1
exec /run.sh
